# include/systemd/warewulfd.service.in
[Unit]
Description=Warewulf cluster management daemon
Documentation=https://warewulf.hpcng.org/
After=network-online.target
AssertFileIsExecutable=@BINDIR@/wwctl

[Service]
Type=simple
User=root
Group=root

ExecStart=@BINDIR@/wwctl server start
ExecReload=/usr/bin/wwctl server reload
ExecStop=/usr/bin/wwctl server stop

PIDFile=/var/run/warewulfd.pid
Restart=always

[Install]
WantedBy=multi-user.target
