# include/systemd/warewulfd.service.in
[Unit]
Description=Warewulf cluster management daemon
Documentation=https://warewulf.org/
After=network-online.target
AssertFileIsExecutable=@BINDIR@/wwctl

[Service]
Type=exec
EnvironmentFile=-/etc/default/warewulfd
User=root
Group=root
ExecStart=@BINDIR@/wwctl server
ExecReload=/bin/kill -HUP "$MAINPID"
Restart=always

[Install]
WantedBy=multi-user.target
